# test/test.service
[Unit]
Description=Test Daemon with watchdog

[Service]
ExecStart=/usr/bin/nodejs src/node-sd-daemon/test/test.js
WatchdogSec=10s
