# test/test-socket.socket
[Unit]
Description=Test-SA Socket

[Socket]
ListenStream=8088
